# templates/systemd/awx-redis.service.j2
[Unit]
Description=AWX Redis ({{ awx_redis_identifier }})
Requires={{ sysd_docker_service }}
After={{ sysd_docker_service }}
DefaultDependencies=no

[Service]
Type=simple
Environment="HOME={{ sysd_docker_unit_home_path }}"
ExecStartPre=-{{ sysd_docker_host_sh }} -c '{{ sysd_docker_host_command }} stop -t {{ sysd_docker_stop_grace }} {{ awx_redis_identifier }} 2>/dev/null || true'
ExecStartPre=-{{ sysd_docker_host_sh }} -c '{{ sysd_docker_host_command }} rm {{ awx_redis_identifier }} 2>/dev/null || true'

ExecStartPre={{ sysd_docker_host_command }} create \
      --rm \
      --name={{ awx_redis_identifier }} \
      --log-driver=none \
      --user=redis:redis \
      --cap-drop=ALL \
      --network={{ awx_redis_container_network }} \
      {{ awx_redis_container_image }}

ExecStart={{ sysd_docker_host_command }} start --attach {{ awx_redis_identifier }}

ExecStop=-{{ sysd_docker_host_sh }} -c '{{ sysd_docker_host_command }} stop -t {{ sysd_docker_stop_grace }} {{ awx_redis_identifier }} 2>/dev/null || true'
ExecStop=-{{ sysd_docker_host_sh }} -c '{{ sysd_docker_host_command }} rm {{ awx_redis_identifier }} 2>/dev/null || true'

Restart=always
RestartSec=30
SyslogIdentifier={{ awx_redis_identifier }}

[Install]
WantedBy=multi-user.target
